# rootfs/etc/supervisor/conf.d/nextcloud-log.conf
[program:nextcloud-log]
command=/usr/local/bin/nextcloud-log
autostart=true
autorestart=true
priority=3
stdout_events_enabled=true
stderr_events_enabled=true
